# alpine-term/services/bb-terminal.ini
[program:bbterminal]
command=/bin/bash /usr/bin/bb-terminal.sh
priority=1
autostart=true
autorestart=true
startsecs=0