# alpine-term/services/program_terminal.conf
[program:bbterminal]
command=/usr/bin/bb-terminal -p 40020 -m 64 --writable -t disableLeaveAlert=true bash
priority=1
autostart=true
autorestart=true
startsecs=0