# alpine-term/services/program_dockerd.conf
[program:dockerd]
command=/usr/bin/dockerd
priority=1
autostart=true
autorestart=true
startsecs=0
